# bin/watchtower.service
[Unit]
Description=Solana Watchtower
Wants=solana-sys-tuner.service
After=network.target

[Service]
User=sol
Type=simple
Restart=always
RestartSec=1
StartLimitInterval=0
Environment="PATH=/home/sol/.local/share/solana/install/active_release/bin"
ExecStart=/bin/bash -c "PATH=/home/sol/bin:$PATH exec /home/solana2/.local/bin/watchtower.sh"

[Install]
WantedBy=multi-user.target
